# data/indicator-location.service.in
[Unit]
Description=Indicator Location Backend
PartOf=graphical-session.target
After=indicators-pre.target

[Service]
ExecStart=@pkglibexecdir@/indicator-location-service
Restart=on-failure
